# Assignment2/testing_review/crontab
*/30 *  *   *   *    /path/to/serveralive.py